# docker/namingo/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/namingo/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:epp]
command=php /opt/registry/epp/start_epp.php
directory=/opt/registry/epp
autostart=true
autorestart=true
stderr_logfile=/var/log/namingo/epp_stderr.log
stdout_logfile=/var/log/namingo/epp_stdout.log
priority=10

[program:whois]
command=php /opt/registry/whois/port43/start_whois.php
directory=/opt/registry/whois/port43
autostart=true
autorestart=true
stderr_logfile=/var/log/namingo/whois_stderr.log
stdout_logfile=/var/log/namingo/whois_stdout.log
priority=20

[program:rdap]
command=php /opt/registry/rdap/start_rdap.php
directory=/opt/registry/rdap
autostart=true
autorestart=true
stderr_logfile=/var/log/namingo/rdap_stderr.log
stdout_logfile=/var/log/namingo/rdap_stdout.log
priority=20

[program:das]
command=php /opt/registry/das/start_das.php
directory=/opt/registry/das
autostart=true
autorestart=true
stderr_logfile=/var/log/namingo/das_stderr.log
stdout_logfile=/var/log/namingo/das_stdout.log
priority=20

[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile
autostart=true
autorestart=true
stderr_logfile=/var/log/namingo/caddy_stderr.log
stdout_logfile=/var/log/namingo/caddy_stdout.log
priority=30

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
stderr_logfile=/var/log/namingo/phpfpm_stderr.log
stdout_logfile=/var/log/namingo/phpfpm_stdout.log
priority=30
